# .config/systemd/user/pacman-notify.timer
[Unit]
Description=Pacman update notifications

[Timer]
OnBootSec=2h
OnUnitActiveSec=5h
RandomizedDelaySec=1h

[Install]
WantedBy=default.target
