# .config/systemd/user/redshift.service
[Unit]
Description=Redshift
StartLimitIntervalSec=0

[Service]
Type=simple
RestartSec=3
Restart=always
SyslogIdentifier=redshift
ExecStart=/usr/bin/env redshift

[Install]
WantedBy=default.target
