# .config/systemd/user/redshift.service
[Unit]
Description=Redshift

[Service]
Type=simple
ExecStart=/usr/bin/env redshift
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target
