# .config/systemd/user/x11vnc.service
[Unit]
Description=Remote desktop service (VNC)
After=syslog.target network.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/env x11vnc -usepw -nowf -repeat -forever -noxrecord -auth %h/.Xauthority
Restart=always
RestartSec=1s

[Install]
WantedBy=default.target
