# .config/systemd/user/x11vnc.service
[Unit]
Description=Remote desktop service (VNC)
StartLimitIntervalSec=0

[Service]
Type=simple
RestartSec=1s
Restart=always
SyslogIdentifier=x11vnc
ExecStart=/usr/bin/env x11vnc -display :0 -rfbport 5900 -shared -usepw -forever -nowireframe -norepeat -noxdamage -noxrecord -ncache_cr -noxfixes -auth $XAUTHORITY

[Install]
WantedBy=default.target
