# .config/systemd/user/x11vnc.service
[Unit]
Description=Remote desktop service (VNC)
StartLimitBurst=5
StartLimitIntervalSec=600

[Service]
Type=simple
RestartSec=3
Restart=always
SyslogIdentifier=x11vnc
ExecStart=/usr/bin/env x11vnc -display $DISPLAY -rfbport 5900 -shared -usepw -forever -nowireframe -norepeat -noxdamage -noxrecord -ncache_cr -noxfixes -auth $XAUTHORITY

[Install]
WantedBy=default.target
