# .config/systemd/user/x11vnc.service
[Unit]
Description=Remote desktop service (VNC)
After=syslog.target network.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/env x11vnc -usepw -nowf -norepeat -forever -noxrecord -auth $XAUTHORITY
Restart=always
RestartSec=1s
SyslogIdentifier=x11vnc

[Install]
WantedBy=default.target
