# .config/systemd/user/xbindkeys.service
[Unit]
Description=Keyboard shortcut binder
StartLimitIntervalSec=0

[Service]
Type=simple
RestartSec=3
Restart=always
KillMode=process
SyslogIdentifier=xbindkeys
ExecStart=/usr/bin/env bash -l -c 'xbindkeys -f %h/.config/X11/xbindkeysrc -p -n'

[Install]
WantedBy=default.target
