# .config/systemd/user/xbindkeys.service
[Unit]
Description=Keyboard shortcut binder

[Service]
Type=simple
RestartSec=5s
Restart=always
KillMode=process
SyslogIdentifier=xbindkeys
EnvironmentFile=-%h/.cache/environment
ExecStart=/usr/bin/env xbindkeys -f %h/.config/xbindkeysrc -p -n
ExecStartPre=/usr/bin/env bash -l -c 'env > %h/.cache/environment'

[Install]
WantedBy=default.target
